# dots/.config/systemd/user/compositor.target
[Unit]
Description=Compositor target for WM to be launched manually
Requires=graphical-session.target
After=graphical-session.target
